# for_systemd/my-web-db-backup.timer
[Unit]
Description=timer for web db backup.sh

[Timer]
OnCalendar=*-*-* *:*:00
Persistent=true

[Install]
WantedBy=timers.target